# deploy/launchagents/ai.renaiss.corerefresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.renaiss.corerefresh</string>
  <key>ProgramArguments</key>
  <array>
    <string>/opt/anaconda3/bin/python3</string>
    <string>/Users/f/Workspace/renaiss_ev_monitor_v2/scripts/refresh_core.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/f/Workspace/renaiss_ev_monitor_v2</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>900</integer>
  <key>StandardOutPath</key>
  <string>/Users/f/.hermes/logs/renaiss-corerefresh.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/f/.hermes/logs/renaiss-corerefresh.stderr.log</string>
</dict>
</plist>
