# deploy/launchagents/ai.renaiss.compare.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Renaiss vs. external price comparison (our independently-scraped public
     market price vs Renaiss FMV). Runs build_comparison.py every 6 hours;
     our_price.py has a 24h cache, so each run only refreshes stale cards and
     does not hammer PriceCharting. Produces data/comparison.json for the
     dashboard /compare page. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.renaiss.compare</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/anaconda3/bin/python3</string>
    <string>/Users/f/Workspace/renaiss_ev_monitor_v2/scripts/build_comparison.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/f/Workspace/renaiss_ev_monitor_v2</string>

  <key>StartInterval</key>
  <integer>21600</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>ThrottleInterval</key>
  <integer>60</integer>

  <key>StandardOutPath</key>
  <string>/Users/f/.hermes/logs/renaiss-compare.stdout.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/f/.hermes/logs/renaiss-compare.stderr.log</string>
</dict>
</plist>
