# deploy/launchagents/ai.renaiss.newpool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- New limited-pack pool capture: scans every 15 minutes. Normally only
     queries the Renaiss API (cheap); only during the open window (from 30 min
     before activeFrom) does it run on-chain probing + pool-wallet analysis.
     Sends a Telegram alert when a new pool appears on-chain. Requires the
     anaconda python (with web3). -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.renaiss.newpool</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/anaconda3/bin/python3</string>
        <string>/Users/f/Workspace/renaiss_ev_monitor_v2/scripts/watch_new_pool.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/f/Workspace/renaiss_ev_monitor_v2</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/anaconda3/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/f/.hermes/logs/renaiss-newpool.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/f/.hermes/logs/renaiss-newpool.err.log</string>
</dict>
</plist>
